# systemd/codex-feishu-daily-stats.timer
[Unit]
Description=Sync Codex Feishu Bridge statistics hourly at minute 10

[Timer]
OnCalendar=*-*-* *:10:00
Persistent=true
AccuracySec=1s
RandomizedDelaySec=0
Unit=codex-feishu-daily-stats.service

[Install]
WantedBy=timers.target
